# etc/s6-overlay/s6-rc.d/init-info/run
#!/command/with-contenv bash
# shellcheck shell=bash

if [[ -n "${REVISION}" ]]; then
  echo "GitRevision: ${REVISION}"
fi

if [[ -n "${VERSION}" ]]; then
  echo "GitVersion: ${VERSION}"
fi

